# packages/engines/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
command=node ./dist/websockets/main.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
